# automaton.service
[Unit]
Description=Slides automation

[Service]
ExecStart=/home/pi/automaton/autostart.sh
Restart=on-failure

[Install]
WantedBy=default.target
